# bizmate-worker.conf
[program:bizmate-worker]
process_name=%(program_name)s_%(process_num)02d
command=/usr/bin/php /var/www/html/unejcat/bizmate.unejcat.com/artisan queue:work --sleep=3 --tries=3 --max-time=3600
autostart=true
autorestart=true
stopasgroup=true
killasgroup=true
user=idtotech
numprocs=1
redirect_stderr=true
stdout_logfile=/var/www/html/unejcat/bizmate.unejcat.com/storage/logs/worker.log
stopwaitsecs=3600
